# startup_scripts/Marcel
#!/sbin/runscript
# Copyright Laurent Faillie 2015
#
# Start/Stop Marcel daemon on OpenRC system (tested on Gentoo)

depend() {
	need net
	after mosquitto
}

TID_BIN=/usr/local/sbin/Marcel
TID_PID=/var/run/Marcel.pid
TID_USR=daemon:daemon

start() {
	ebegin "Starting Marcel daemon"
	start-stop-daemon --start --exec ${TID_BIN} \
		--interpreted --user ${TID_USR} \
		--background --make-pidfile --pidfile ${TID_PID}
	eend $?
}

stop() {
	ebegin "Stopping Marcel daemon"
	start-stop-daemon --stop --exec ${TID_BIN} \
		--pidfile ${TID_PID}
	eend $?
}

